# rules/atomics/T1547.008/T1547_008_Modify_Registry_to_load_Arbitrary_DLL_into_LSASS_-.toml
[rule]
author = [ "AtomicRedTeam-Automation",]
description = "Detection for specific test 'Modify Registry to load Arbitrary DLL into LSASS - LsaDbExtPt' (GUID: 8ecef16d-d289-46b4-917b-0dba6dc81cf1) of technique T1547.008 (Boot or Logon Autostart Execution: LSASS Driver). Test Description: The following Atomic will modify an undocumented registry key that may be abused to load a arbitrary DLL into LSASS. \n\nUpon execution, the registry key will be modified and a value will contain the path to the DLL. \nReference: https://blog.xpnsec.com/exploring-mimikatz-part-1/ and source https://github.com/oxfemale/LogonCredentialsSteal\nNote that if any LSA based protection is enabled, this will most likely not be successful with LSASS.exe loading the DLL.. Command/Steps: New-ItemProperty -Path HKLM:\\\\SYSTEM\\\\CurrentControlSet\\\\Services\\\\NTDS -Name LsaDbExtPt -Value \\\"#{dll_path}\\\"\n"
from = "now-9m"
rule_id = "518a1900-9162-43c5-9152-7b0d2720ea9c"
language = "kql"
name = "Atomic Test: T1547.008 - Modify Registry to load Arbitrary DLL into LSASS - LsaDbExtPt"
output_index = ".siem-signals-default"
references = [ "https://raw.githubusercontent.com/redcanaryco/atomic-red-team/refs/heads/master/atomics/T1547.008/T1547.008.yaml",]
risk_score = 50
severity = "medium"
tags = [ "AtomicRedTeam", "T1547.008", "8ecef16d-d289-46b4-917b-0dba6dc81cf1", "windows",]
type = "query"
query = "(process.name : (\\\"powershell.exe\\\", \\\"pwsh.exe\\\", \\\"powershell_ise.exe\\\") OR process.parent.name : (\\\"powershell.exe\\\", \\\"pwsh.exe\\\", \\\"powershell_ise.exe\\\")) AND (powershell.script_block_text : \\\"*New-ItemProperty*\\\" AND powershell.script_block_text : \\\"*-Path*\\\" AND powershell.script_block_text : \\\"*HKLM:\\\\SYSTEM\\\\CurrentControlSet\\\\Services\\\\NTDS*\\\" OR (process.command_line : \\\"*New-ItemProperty*\\\" AND process.command_line : \\\"*-Path*\\\" AND process.command_line : \\\"*HKLM:\\\\SYSTEM\\\\CurrentControlSet\\\\Services\\\\NTDS*\\\"))"
version = 1
index = [ "logs-endpoint.events.*", "winlogbeat-*", "logs-windows.*", "logs-linux.*", "logs-macos.*",]
[[rule.threat]]
framework = "MITRE ATT&CK"
[[rule.threat.technique]]
id = "T1547.008"
name = "Modify Registry to load Arbitrary DLL into LSASS - LsaDbExtPt"
reference = "https://attack.mitre.org/techniques/T1547/008"

[[rule.threat.tactic]]
id = "TA0003"
name = "Persistence"
reference = "https://attack.mitre.org/techniques/T1547/008"

[[rule.threat.tactic]]
id = "TA0004"
name = "Privilege Escalation"
reference = "https://attack.mitre.org/techniques/T1547/008"

[[rule.threat.tactic]]
id = "TA0005"
name = "Defense Evasion"
reference = "https://attack.mitre.org/techniques/T1547/008"


