# rules/atomics/T1546/T1546_Persistence_with_Custom_AutodialDLL.toml
[rule]
author = [ "AtomicRedTeam-Automation",]
description = "Detection for specific test 'Persistence with Custom AutodialDLL' (GUID: aca9ae16-7425-4b6d-8c30-cad306fdbd5b) of technique T1546 (Event Triggered Execution). Test Description: The DLL pointed to by the AutodialDLL registry key is loaded every time a process connects to the internet. Attackers can gain persistent code execution by setting this key to a DLL of their choice. \n\nThe sample dll provided, AltWinSock2DLL, will launch the notepad process. Starting and stopping a web browser such as MS Edge or Chrome should result in the dll executing.\n[Blog](https://www.mdsec.co.uk/2022/10/autodialdlling-your-way/). Command/Steps: Set-ItemProperty HKLM:\\\\SYSTEM\\\\CurrentControlSet\\\\Services\\\\WinSock2\\\\Parameters -Name AutodialDLL -Value PathToAtomicsFolder\\\\T1546\\\\bin\\\\AltWinSock2DLL.dll..."
from = "now-9m"
rule_id = "34b19c9c-8263-45df-9568-954e81229807"
language = "kql"
name = "Atomic Test: T1546 - Persistence with Custom AutodialDLL"
output_index = ".siem-signals-default"
references = [ "https://raw.githubusercontent.com/redcanaryco/atomic-red-team/refs/heads/master/atomics/T1546/T1546.yaml",]
risk_score = 50
severity = "medium"
tags = [ "AtomicRedTeam", "T1546", "aca9ae16-7425-4b6d-8c30-cad306fdbd5b", "windows",]
type = "query"
query = "(process.name : (\\\"powershell.exe\\\", \\\"pwsh.exe\\\", \\\"powershell_ise.exe\\\") OR process.parent.name : (\\\"powershell.exe\\\", \\\"pwsh.exe\\\", \\\"powershell_ise.exe\\\")) AND (powershell.script_block_text : \\\"*Set-ItemProperty*\\\" AND powershell.script_block_text : \\\"*HKLM:\\\\SYSTEM\\\\CurrentControlSet\\\\Services\\\\WinSock2\\\\Parameters*\\\" AND powershell.script_block_text : \\\"*-Name*\\\" OR (process.command_line : \\\"*Set-ItemProperty*\\\" AND process.command_line : \\\"*HKLM:\\\\SYSTEM\\\\CurrentControlSet\\\\Services\\\\WinSock2\\\\Parameters*\\\" AND process.command_line : \\\"*-Name*\\\"))"
version = 1
index = [ "logs-endpoint.events.*", "winlogbeat-*", "logs-windows.*", "logs-linux.*", "logs-macos.*",]
[[rule.threat]]
framework = "MITRE ATT&CK"
[[rule.threat.technique]]
id = "T1546"
name = "Persistence with Custom AutodialDLL"
reference = "https://attack.mitre.org/techniques/T1546"

[[rule.threat.tactic]]
id = "TA0003"
name = "Persistence"
reference = "https://attack.mitre.org/techniques/T1546"

[[rule.threat.tactic]]
id = "TA0004"
name = "Privilege Escalation"
reference = "https://attack.mitre.org/techniques/T1546"

[[rule.threat.tactic]]
id = "TA0005"
name = "Defense Evasion"
reference = "https://attack.mitre.org/techniques/T1546"


