# rules/atomics/T1053.003/T1053_003_Cron_-_Add_script_to__etc_cron_d_folder.toml
[rule]
author = [ "AtomicRedTeam-Automation",]
description = "Detection for specific test 'Cron - Add script to /etc/cron.d folder' (GUID: 078e69eb-d9fb-450e-b9d0-2e118217c846) of technique T1053.003 (Scheduled Task/Job: Cron). Test Description: This test adds a script to /etc/cron.d folder configured to execute on a schedule.. Command/Steps: echo \\\"#{command}\\\" > /etc/cron.d/#{cron_script_name}\n"
from = "now-9m"
rule_id = "acc1fd8b-9491-415d-aad9-b2f84683aeee"
language = "kql"
name = "Atomic Test: T1053.003 - Cron - Add script to /etc/cron.d folder"
output_index = ".siem-signals-default"
references = [ "https://raw.githubusercontent.com/redcanaryco/atomic-red-team/refs/heads/master/atomics/T1053.003/T1053.003.yaml",]
risk_score = 50
severity = "medium"
tags = [ "AtomicRedTeam", "T1053.003", "078e69eb-d9fb-450e-b9d0-2e118217c846", "linux",]
type = "query"
query = "(process.name : (\\\"sh\\\", \\\"bash\\\") OR process.parent.name : (\\\"sh\\\", \\\"bash\\\")) AND process.command_line : \\\"*root*\\\" AND process.command_line : \\\"*Hello*\\\" AND process.command_line : \\\"*from*\\\""
version = 1
index = [ "logs-endpoint.events.*", "winlogbeat-*", "logs-windows.*", "logs-linux.*", "logs-macos.*",]
[[rule.threat]]
framework = "MITRE ATT&CK"
[[rule.threat.technique]]
id = "T1053.003"
name = "Cron - Add script to /etc/cron.d folder"
reference = "https://attack.mitre.org/techniques/T1053/003"

[[rule.threat.tactic]]
id = "TA0003"
name = "Persistence"
reference = "https://attack.mitre.org/techniques/T1053/003"

[[rule.threat.tactic]]
id = "TA0004"
name = "Privilege Escalation"
reference = "https://attack.mitre.org/techniques/T1053/003"

[[rule.threat.tactic]]
id = "TA0005"
name = "Defense Evasion"
reference = "https://attack.mitre.org/techniques/T1053/003"


