# rules/atomics/T1547.001/T1547_001_Suspicious_vbs_file_run_from_startup_Folder.toml
[rule]
author = [ "AtomicRedTeam-Automation",]
description = "Detection for specific test 'Suspicious vbs file run from startup Folder' (GUID: 2cb98256-625e-4da9-9d44-f2e5f90b8bd5) of technique T1547.001 (Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder). Test Description: vbs files can be placed in and ran from the startup folder to maintain persistance. Upon execution, \\\"T1547.001 Hello, World VBS!\\\" will be displayed twice. \nAdditionally, the new files can be viewed in the \\\"$env:APPDATA\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup\\\"\nfolder and will also run when the computer is restarted and the user logs in.. Command/Steps: Copy-Item \\\"$PathToAtomicsFolder\\\\T1547.001\\\\src\\\\vbsstartup.vbs\\\" \\\"$env:APPDATA\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup\\\\vbsstartup.vbs\\\"\nCopy-Item \\\"$..."
from = "now-9m"
rule_id = "159d99a4-d827-4c78-8a2f-bf25da8bd7bf"
language = "kql"
name = "Atomic Test: T1547.001 - Suspicious vbs file run from startup Folder"
output_index = ".siem-signals-default"
references = [ "https://raw.githubusercontent.com/redcanaryco/atomic-red-team/refs/heads/master/atomics/T1547.001/T1547.001.yaml",]
risk_score = 50
severity = "medium"
tags = [ "AtomicRedTeam", "T1547.001", "2cb98256-625e-4da9-9d44-f2e5f90b8bd5", "windows",]
type = "query"
query = "(process.name : (\\\"powershell.exe\\\", \\\"pwsh.exe\\\", \\\"powershell_ise.exe\\\") OR process.parent.name : (\\\"powershell.exe\\\", \\\"pwsh.exe\\\", \\\"powershell_ise.exe\\\")) AND (powershell.script_block_text : \\\"*Copy-Item*\\\" AND powershell.script_block_text : \\\"*$env:APPDATA\\\\Microsoft\\\\Windows\\\\Start*\\\" AND powershell.script_block_text : \\\"*Menu\\\\Programs\\\\Startup\\\\vbsstartup.vbs*\\\" OR (process.command_line : \\\"*Copy-Item*\\\" AND process.command_line : \\\"*$env:APPDATA\\\\Microsoft\\\\Windows\\\\Start*\\\" AND process.command_line : \\\"*Menu\\\\Programs\\\\Startup\\\\vbsstartup.vbs*\\\"))"
version = 1
index = [ "logs-endpoint.events.*", "winlogbeat-*", "logs-windows.*", "logs-linux.*", "logs-macos.*",]
[[rule.threat]]
framework = "MITRE ATT&CK"
[[rule.threat.technique]]
id = "T1547.001"
name = "Suspicious vbs file run from startup Folder"
reference = "https://attack.mitre.org/techniques/T1547/001"

[[rule.threat.tactic]]
id = "TA0003"
name = "Persistence"
reference = "https://attack.mitre.org/techniques/T1547/001"

[[rule.threat.tactic]]
id = "TA0004"
name = "Privilege Escalation"
reference = "https://attack.mitre.org/techniques/T1547/001"

[[rule.threat.tactic]]
id = "TA0005"
name = "Defense Evasion"
reference = "https://attack.mitre.org/techniques/T1547/001"


