# rules/atomics/T1547.001/T1547_001_Suspicious_jse_file_run_from_startup_Folder.toml
[rule]
author = [ "AtomicRedTeam-Automation",]
description = "Detection for specific test 'Suspicious jse file run from startup Folder' (GUID: dade9447-791e-4c8f-b04b-3a35855dfa06) of technique T1547.001 (Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder). Test Description: jse files can be placed in and ran from the startup folder to maintain persistance.\nUpon execution, \\\"T1547.001 Hello, World JSE!\\\" will be displayed twice. \nAdditionally, the new files can be viewed in the \\\"$env:APPDATA\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup\\\"\nfolder and will also run when the computer is restarted and the user logs in.. Command/Steps: Copy-Item \\\"$PathToAtomicsFolder\\\\T1547.001\\\\src\\\\jsestartup.jse\\\" \\\"$env:APPDATA\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup\\\\jsestartup.jse\\\"\nCopy-Item \\\"$..."
from = "now-9m"
rule_id = "bb94c1fb-a0c6-4989-82c1-d2dca9e22ac6"
language = "kql"
name = "Atomic Test: T1547.001 - Suspicious jse file run from startup Folder"
output_index = ".siem-signals-default"
references = [ "https://raw.githubusercontent.com/redcanaryco/atomic-red-team/refs/heads/master/atomics/T1547.001/T1547.001.yaml",]
risk_score = 50
severity = "medium"
tags = [ "AtomicRedTeam", "T1547.001", "dade9447-791e-4c8f-b04b-3a35855dfa06", "windows",]
type = "query"
query = "(process.name : (\\\"powershell.exe\\\", \\\"pwsh.exe\\\", \\\"powershell_ise.exe\\\") OR process.parent.name : (\\\"powershell.exe\\\", \\\"pwsh.exe\\\", \\\"powershell_ise.exe\\\")) AND (powershell.script_block_text : \\\"*Copy-Item*\\\" AND powershell.script_block_text : \\\"*$env:APPDATA\\\\Microsoft\\\\Windows\\\\Start*\\\" AND powershell.script_block_text : \\\"*Menu\\\\Programs\\\\Startup\\\\jsestartup.jse*\\\" OR (process.command_line : \\\"*Copy-Item*\\\" AND process.command_line : \\\"*$env:APPDATA\\\\Microsoft\\\\Windows\\\\Start*\\\" AND process.command_line : \\\"*Menu\\\\Programs\\\\Startup\\\\jsestartup.jse*\\\"))"
version = 1
index = [ "logs-endpoint.events.*", "winlogbeat-*", "logs-windows.*", "logs-linux.*", "logs-macos.*",]
[[rule.threat]]
framework = "MITRE ATT&CK"
[[rule.threat.technique]]
id = "T1547.001"
name = "Suspicious jse file run from startup Folder"
reference = "https://attack.mitre.org/techniques/T1547/001"

[[rule.threat.tactic]]
id = "TA0003"
name = "Persistence"
reference = "https://attack.mitre.org/techniques/T1547/001"

[[rule.threat.tactic]]
id = "TA0004"
name = "Privilege Escalation"
reference = "https://attack.mitre.org/techniques/T1547/001"

[[rule.threat.tactic]]
id = "TA0005"
name = "Defense Evasion"
reference = "https://attack.mitre.org/techniques/T1547/001"


