# rules/atomics/T1543.001/T1543_001_Launch_Agent_-_Root_Directory.toml
[rule]
author = [ "AtomicRedTeam-Automation",]
description = "Detection for specific test 'Launch Agent - Root Directory' (GUID: 66774fa8-c562-4bae-a58d-5264a0dd9dd7) of technique T1543.001 (Create or Modify System Process: Launch Agent). Test Description: Create a plist and execute it. Command/Steps: sudo cp #{path_malicious_plist} /Library/LaunchAgents/#{plist_filename}\nlaunchctl load -w /Library/LaunchAgents/#{plist_filename}\n"
from = "now-9m"
rule_id = "0778fdb7-5a43-42cc-9ce2-117f6aa6afb3"
language = "kql"
name = "Atomic Test: T1543.001 - Launch Agent - Root Directory"
output_index = ".siem-signals-default"
references = [ "https://raw.githubusercontent.com/redcanaryco/atomic-red-team/refs/heads/master/atomics/T1543.001/T1543.001.yaml",]
risk_score = 50
severity = "medium"
tags = [ "AtomicRedTeam", "T1543.001", "66774fa8-c562-4bae-a58d-5264a0dd9dd7", "macos",]
type = "query"
query = "(process.name : (\\\"sh\\\", \\\"bash\\\") OR process.parent.name : (\\\"sh\\\", \\\"bash\\\")) AND process.command_line : \\\"*launchctl*\\\" AND process.command_line : \\\"*load*\\\""
version = 1
index = [ "logs-endpoint.events.*", "winlogbeat-*", "logs-windows.*", "logs-linux.*", "logs-macos.*",]
[[rule.threat]]
framework = "MITRE ATT&CK"
[[rule.threat.technique]]
id = "T1543.001"
name = "Launch Agent - Root Directory"
reference = "https://attack.mitre.org/techniques/T1543/001"

[[rule.threat.tactic]]
id = "TA0003"
name = "Persistence"
reference = "https://attack.mitre.org/techniques/T1543/001"

[[rule.threat.tactic]]
id = "TA0004"
name = "Privilege Escalation"
reference = "https://attack.mitre.org/techniques/T1543/001"

[[rule.threat.tactic]]
id = "TA0005"
name = "Defense Evasion"
reference = "https://attack.mitre.org/techniques/T1543/001"


