# rules/atomics/T1053.005/T1053_005_WMI_Invoke-CimMethod_Scheduled_Task.toml
[rule]
author = [ "AtomicRedTeam-Automation",]
description = "Detection for specific test 'WMI Invoke-CimMethod Scheduled Task' (GUID: e16b3b75-dc9e-4cde-a23d-dfa2d0507b3b) of technique T1053.005 (Scheduled Task/Job: Scheduled Task). Test Description: Create an scheduled task that executes notepad.exe after user login from XML by leveraging WMI class PS_ScheduledTask. Does the same thing as Register-ScheduledTask cmdlet behind the scenes.. Command/Steps: $xml = [System.IO.File]::ReadAllText(\\\"#{xml_path}\\\")\nInvoke-CimMethod -ClassName PS_ScheduledTask -NameSpace \\\"Root\\\\Microsoft\\\\Windows\\\\TaskScheduler\\\" -Me..."
from = "now-9m"
rule_id = "3f79a6e7-cb7f-4835-8eda-36f4b4d2fb9e"
language = "kql"
name = "Atomic Test: T1053.005 - WMI Invoke-CimMethod Scheduled Task"
output_index = ".siem-signals-default"
references = [ "https://raw.githubusercontent.com/redcanaryco/atomic-red-team/refs/heads/master/atomics/T1053.005/T1053.005.yaml",]
risk_score = 50
severity = "medium"
tags = [ "AtomicRedTeam", "T1053.005", "e16b3b75-dc9e-4cde-a23d-dfa2d0507b3b", "windows",]
type = "query"
query = "(process.name : (\\\"powershell.exe\\\", \\\"pwsh.exe\\\", \\\"powershell_ise.exe\\\") OR process.parent.name : (\\\"powershell.exe\\\", \\\"pwsh.exe\\\", \\\"powershell_ise.exe\\\")) AND (powershell.script_block_text : \\\"*$xml*\\\" AND powershell.script_block_text : \\\"*Invoke-CimMethod*\\\" AND powershell.script_block_text : \\\"*-ClassName*\\\" OR (process.command_line : \\\"*$xml*\\\" AND process.command_line : \\\"*Invoke-CimMethod*\\\" AND process.command_line : \\\"*-ClassName*\\\"))"
version = 1
index = [ "logs-endpoint.events.*", "winlogbeat-*", "logs-windows.*", "logs-linux.*", "logs-macos.*",]
[[rule.threat]]
framework = "MITRE ATT&CK"
[[rule.threat.technique]]
id = "T1053.005"
name = "WMI Invoke-CimMethod Scheduled Task"
reference = "https://attack.mitre.org/techniques/T1053/005"

[[rule.threat.tactic]]
id = "TA0003"
name = "Persistence"
reference = "https://attack.mitre.org/techniques/T1053/005"

[[rule.threat.tactic]]
id = "TA0004"
name = "Privilege Escalation"
reference = "https://attack.mitre.org/techniques/T1053/005"

[[rule.threat.tactic]]
id = "TA0005"
name = "Defense Evasion"
reference = "https://attack.mitre.org/techniques/T1053/005"


